# fanspeed.service
[Unit]
Description=

[Service]
Type=exec
ExecStart=/usr/local/bin/fanspeed
